# 50-eth-offload.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="enp0s20f0u3", RUN+="/usr/bin/ethtool -K enp0s20f0u3 gro on lro on tso on"
